# gunicorn.conf
[program:gunicorn]
command=/usr/bin/gunicorn WordsFromHeadlines.wsgi -b localhost:8000
directory=/deploy/app
